# systemd/nucleus.service
[Unit]
Description=tokens distribution service

[Service]
ExecStart=/usr/bin/nucleus
Restart=always

[Install]
WantedBy=multi-user.target
